# deploy/systemd/nginx.service.d/hardening.conf
# =============================================================================
# Nginx service hardening drop-in.
# Installed to /etc/systemd/system/nginx.service.d/hardening.conf by
# scripts/install-target.sh. Complements the limits.conf drop-in that
# apply-layer-1.sh writes (LimitNOFILE) — keep the two separate.
#
# Nginx is special: the master starts as root to bind :80/:443 and to setuid
# workers down to the `nginx` user, so it keeps exactly three capabilities and
# nothing else. Workers themselves run unprivileged.
# =============================================================================
[Service]
# Only the caps nginx genuinely needs: bind privileged ports + drop to `nginx`.
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_NET_BIND_SERVICE

NoNewPrivileges=false          # master must setuid workers; cannot be true here
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
# Nginx needs to write logs, the pid file, and its temp/cache dirs.
ReadWritePaths=/var/log/nginx /var/lib/nginx /var/cache/nginx /run
